# anglia/dgv0.crontab
#m      h       dom     mon     dow     command
30 * * * * bash /home/dgv0/scripts/ip-checker.sh
