# anglia/minecraft.crontab
#m      h       dom     mon     dow     command
0 6 * * * /usr/sbin/service minecraft restart >>/srv/log/cron.log 2>&1
0 8 * * * /usr/sbin/service minecraft backup >>/srv/log/cron.log 2>&1
